# sv/session/run
#!/bin/sh
cgcreate -g pids:/juan
exec cgexec -g pids:/juan --sticky -- su juan -l -c ~juan/env/init
