# xinit.d/services/alacritty/run
#!/bin/sh
exec alacritty --config-file ./alacritty.yml --working-directory ~ -e tmux -2 new-session -A -s main
