# xinit.d/services/alacritty/run
#!/bin/sh
exec alacritty --config-file ./alacritty.yml -e tmux attach-session -t main
#--working-directory ~
#-e tmux -2 new-session -A -s main
