# xinit.d/services/kitty/run
#!/bin/sh
_cwd=$(pwd); cd ~
exec kitty --config $_cwd/kitty.conf 2>&1
