# sv/wezterm/run
#!/bin/bash
exec env -C ~ wezterm --config-file $(dirname $(realpath $0))/wezterm.lua start --class terminal 2>&1
