# sv/wezterm/run
#!/bin/sh
_cfg=$(pwd)/wezterm.lua
exec env -C ~ wezterm --config-file $_cfg start --class terminal 2>&1
