# services/dnscrypt-proxy/run
#!/bin/sh
exec dnscrypt-proxy -config ./dnscrypt-proxy.toml
