# services/dnscrypt-proxy/run
#!/bin/sh
exec dnscrypt-proxy -config ~juan/env/dns.d/dnscrypt-proxy.toml
