# xinit.d/services/chromium/run
#!/bin/sh
cd ~
NO_AT_BRIDGE=1 exec chromium 2>&1
