# services/connmand/run
#!/bin/sh
exec connmand --nodaemon --nodnsproxy --noplugin=vpn --config=./connmand.conf
