# services/connmand/run
#!/bin/sh
exec connmand --nodaemon --nodnsproxy --config=./connmand.conf
