# services/wg/run
#!/bin/sh
wg-quick up ./hx.conf || exit 1
exec chpst -b wg pause
