# sv/dunst/run
#!/bin/sh
exec dunst -config ./config 2>&1
