# services/caddy/run
#!/bin/sh
setcap cap_net_bind_service=+ep /usr/bin/caddy
CADDYPATH=. exec chpst -o 8192 -u juan caddy