# services/caddy/run
#!/bin/sh
exec chpst -o 8192 -u juan caddy
