# xinit.d/services/pipewire-pulse/run
#!/bin/sh
vx check pipewire
exec pipewire-pulse 2>&1
