# xinit.d/services/sxhkd/run
#!/bin/sh
_cwd=$(pwd); cd ~
exec sxhkd -c $_cwd/config 2>&1
