# xinit.d/services/sxhkd/run
#!/bin/sh
_cfg=$(pwd)/config
exec env -C ~ sxhkd -c $_cfg 2>&1
