# xinit.d/services/wireplumber/run
#!/bin/sh
vx check pipewire >/dev/null || exit 1
exec wireplumber 2>&1
