# services/autologin/run
#!/bin/sh
exec setsid agetty --noclear --autologin juan tty1 linux